# crontab
# DO NOT EDIT THIS FILE - edit the master and reinstall.
# (/tmp/crontab.Znp648/crontab installed on Thu Jun 27 23:56:01 2024)
# (Cron version -- $Id: crontab.c,v 2.13 1994/01/17 03:20:37 vixie Exp $)
# Edit this file to introduce tasks to be run by cron.
# 
# Each task to run has to be defined through a single line
# indicating with different fields when the task will be run
# and what command to run for the task
# 
# To define the time you can provide concrete values for
# minute (m), hour (h), day of month (dom), month (mon),
# and day of week (dow) or use '*' in these fields (for 'any').
# 
# Notice that tasks will be started based on the cron's system
# daemon's notion of time and timezones.
# 
# Output of the crontab jobs (including errors) is sent through
# email to the user the crontab file belongs to (unless redirected).
# 
# For example, you can run a backup of all your user accounts
# at 5 a.m every week with:
# 0 5 * * 1 tar -zcf /var/backups/home.tgz /home/
# 
# For more information see the manual pages of crontab(5) and cron(8)
# 
# m h  dom mon dow   command


###Caches Clenaer

*/30 * * * * sudo sh -c 'echo 1 > /proc/sys/vm/drop_caches && echo "$(date): \e[37mCache level\e[0m \e[31m2\e[0m \e[32mdropped successfully\e[0m" >> /var/log/cleaner/Caches-Cleaner/drop_caches_1.log'
*/30 * * * * sudo sh -c 'echo 2 > /proc/sys/vm/drop_caches && echo "$(date): \e[37mCache level\e[0m \e[31m1\e[0m \e[32mdropped successfully\e[0m" >> /var/log/cleaner/Caches-Cleaner/drop_caches_1.log'
*/30 * * * sudo sh -c 'echo 3 > /proc/sys/vm/drop_caches && echo "$(date): \e[37mCache level\e[0m \e[31m3\e[0m \e[32mdropped successfully\e[0m" >> /var/log/cleaner/Caches-Cleaner/drop_caches_1.log'

*/20 * * * * sudo sh -c 'sync; echo 1 > /proc/sys/vm/drop_caches && echo "$(date): \e[37mCache level\e[0m \e[31m1\e[0m \e[32mdropped with sync successfully\e[0m" >> /var/log/cleaner/Caches-Cleaner/drop_caches_2.log'
*/20 * * * * sudo sh -c 'sync; echo 2 > /proc/sys/vm/drop_caches && echo "$(date): \e[37mCache level\e[0m \e[31m3\e[0m \e[32mdropped with sync successfully\e[0m" >> /var/log/cleaner/Caches-Cleaner/drop_caches_2.log'
*/20 * * * * sudo sh -c 'sync; echo 3 > /proc/sys/vm/drop_caches && echo "$(date): \e[37mCache level\e[0m \e[31m2\e[0m \e[32mdropped with sync successfully\e[0m" >> /var/log/cleaner/Caches-Cleaner/drop_caches_2.log'




###Clean the drop_caches_1.log
*/31 * * * * truncate -s 0 /var/log/cleaner/Caches-Cleaner/drop_caches_1.log
*/21 * * * * truncate -s 0 /var/log/cleaner/Caches-Cleaner/drop_caches_2.log



###System Cleaner

#0 2 * * 0 journalctl --vacuum-time=1w | sudo tee -a /var/log/cleaner/cron_journalctl.log > /dev/null
0 2 * * 0 journalctl --vacuum-time=1w && echo "\e[32m   Successfully cleaned\e[0m journalctl logs at $(date)" | sudo tee -a /var/log/cleaner/Systems-Cleaner/cron_journalctl.log > /dev/null

#0 3 * * 0 rm -rf /tmp/* | sudo tee -a /var/log/cleaner/cron_tmp_clean.log > /dev/null
0 3 * * 0 rm -rf /tmp/* && echo "\e[32mSuccessfully cleaned\e[0m /tmp directory at $(date)" | sudo tee -a /var/log/cleaner/Systems-Cleaner/cron_tmp_clean.log > /dev/null

#0 4 * * 0 rm -rf ~/.cache/* | sudo tee -a /var/log/cleaner/cron_cache_clean.log > /dev/null
0 4 * * 0 rm -rf ~/.cache/* && echo "\e[32mSuccessfully cleaned\e[0m ~/.cache directory at $(date)" | sudo tee -a /var/log/cleaner/Systems-Cleaner/cron_cache_clean.log > /dev/null


###Clean the System Cleaner every Monday
0 5 * * 0 truncate -s 0 /var/log/cleaner/Systems-Cleaner/cron_journalctl.log
0 5 * * 0 truncate -s 0 /var/log/cleaner/Systems-Cleaner/cron_tmp_clean.log
0 5 * * 0 truncate -s 0 /var/log/cleaner/Systems-Cleaner/cron_cache_clean.log


#* * * * * echo "Test log entry for journalctl vacuum" | sudo tee -a /var/log/cleaner/cron_journalctl.log > /dev/null
#* * * * * echo "Test log entry for /tmp clean" | sudo tee -a /var/log/cleaner/cron_tmp_clean.log > /dev/null
#* * * * * echo "Test log entry for ~/.cache clean" | sudo tee -a /var/log/cleaner/cron_cache_clean.log > /dev/null





###Apt

# Run apt-get clean --dry-run at 1 AM daily
0 1 * * * sudo apt-get clean --dry-run && echo "$(date): \e[32mSuccessfully\e[0m ran 'apt-get clean --dry-run'" >> /var/log/cleaner/Apt-Cleaner/clean_apt.log

# Run apt-get autoclean at 2 AM daily
0 2 * * * sudo apt-get autoclean && echo "$(date): \e[32mSuccessfully\e[0m ran 'apt-get autoclean'" >> /var/log/cleaner/Apt-Cleaner/clean_apt.log

# Run apt-get autoremove at 3 AM daily
0 3 * * * sudo apt-get autoremove --purge -y && echo "$(date): \e[32mSuccessfully\e[0m ran 'apt-get autoremove'" >> /var/log/cleaner/Apt-Cleaner/clean_apt.log

###Clean Apt Logs every 2 days
0 4 * * * truncate -s 0 /var/log/cleaner/Apt-Cleaner/clean_apt.log


###Processe
*/30 * * * * /var/log/cleaner/Processes/cleanup.sh

*/31 * * * * truncate -s 0 /var/log/cleaner/Processes/process_cleanup.log

###Update The Systems
0 */6 * * * /var/log/cleaner/Update-Repository/apt_update.sh




#* * * * * /var/log/cleaner/Troubleshoot/collect_troubleshoot_logs.sh >/dev/null 2>&1

#* * * * * /var/log/cleaner/Bell/play_sound.sh && /var/log/cleaner/Update-Repository/apt_update.sh










